# src/hdp/crontab
#
# m h  dom mon dow   command

00 09  * * * sudo bash /home/xxxxx/test/00_330.sh